# Packer/test-backend/setup-metering.service
[Unit]
Description=Setup metering reporting

[Service]
Type=oneshot
ExecStart=/app/setup-metering.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
